# etc/startup_mediaserver.conf
[program:mediaserver]
command=/opt/janus/bin/janus
autostart=true
autorestart=true
stderr_logfile=/var/log/mediaserver.log
stdout_logfile=/var/log/mediaserver.log